# etc/crontab/rca.crontab
# CRONTAB rca:exptools
#
#
# min	hour	daymo	month	daywk	command
#
#	min  - time(s) of day 
#	hour
#	daymo - day(s) of month (1, 2, ... 31)
#	month - month(s) of the year (1, 2, ... 12)
#	daywk - day(s) of week (0-6, 0 = sun)
#
#
# min	hour	daymo	month	daywk	command
#
#
# GROPE maintenance
#
7 8-23 * * * sbin/gropedict.cron
#
#
# permissions (mostly make readable)
#
10 2 * * * sbin/perms.cron
#
#
#
